# ops/systemd/matrixcode.service
[Unit]
Description=MatrixCode intelligent agent console
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=matrixcode
Group=matrixcode
WorkingDirectory=/opt/matrixcode
EnvironmentFile=/etc/matrixcode/matrixcode.env
Environment=MATRIXCODE_APP_HOME=/opt/matrixcode
Environment=MATRIXCODE_ENV_FILE=/etc/matrixcode/matrixcode.env
ExecStart=/opt/matrixcode/bin/run-matrixcode-server.sh
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/matrixcode /var/log/matrixcode

[Install]
WantedBy=multi-user.target
